# platform/broadcom/sonic-platform-modules-cel/debian/platform-modules-ds1000.postinst
depmod -a
/usr/local/bin/platform_api_mgnt.sh install
systemctl enable pddf-platform-init.service
systemctl enable ds1000-fan-control.service
systemctl start pddf-platform-init.service
systemctl start ds1000-fan-control.service
