# platform/broadcom/sonic-platform-modules-cel/debian/platform-modules-questone2.postinst
depmod -a

/usr/local/bin/platform_api_mgnt.sh install

systemctl enable platform-modules-questone2.service
systemctl start platform-modules-questone2.service
