# platform/broadcom/sonic-platform-modules-nexthop/common/etc/systemd/system/syncd.service.d/override.conf
[Unit]
Wants=pmon.service
After=pmon.service
