# platform/broadcom/sonic-platform-modules-cel/ds1000/service/ds1000-fan-control.service
[Unit]
Description=Ds1000 Fan Control service
After=pddf-platform-init.service
Requires=pddf-platform-init.service
BindsTo=pddf-platform-init.service

[Service]
Type=simple
ExecStart=/usr/local/bin/ds1000_fanctld.py
KillSignal=SIGTERM
SuccessExitStatus=0

# Resource Limitations
LimitCORE=infinity

[Install]
WantedBy=multi-user.target
