# platform/broadcom/sonic-platform-modules-nexthop/common/etc/systemd/system/determine-reboot-cause.service.d/override.conf
[Unit]
Requires=pddf-platform-init.service
After=pddf-platform-init.service
