# platform/broadcom/sonic-platform-modules-cel/silverstone-v2/service/silverstone-v2-pddf-platform-monitor.service
[Unit]
Description=Silverstone-V2 Platform background monitor service
Before=pmon.service
After=pddf-platform-init.service
Requires=pddf-platform-init.service
DefaultDependencies=no

[Service]
ExecStart=/usr/local/bin/pddf_fan_control_sensor_refresh.py
Restart=always
RestartSec=30s
KillSignal=SIGKILL
SuccessExitStatus=SIGKILL

# Resource Limitations
LimitCORE=infinity

[Install]
WantedBy=multi-user.target
