# platform/broadcom/sonic-platform-modules-dell/debian/platform-modules-z9664f.prerm
# prerm script for Z9664f

systemctl stop platform-modules-z9664f.service
systemctl disable platform-modules-z9664f.service

