# platform/broadcom/sonic-platform-modules-nexthop/common/service/adm1266-rtc-sync.service
[Unit]
Description=ADM1266 RTC synchronization service
After=pddf-platform-init.service
Requires=pddf-platform-init.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/adm1266_rtc_sync

[Install]
WantedBy=multi-user.target
