# platform/broadcom/sonic-platform-modules-nexthop/nh-4010/etc/systemd/system/syncd.service.d/ge-override.conf
[Service]
ExecStopPost=/usr/local/bin/asic_init.sh
