# platform/broadcom/sonic-platform-modules-micas/m2-w6940-128qc/service/subnetwork.service
[Unit]
Description=Configuration eth0.4088 subnetwork
After=network.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/local/bin/subnetwork.py start
ExecStop=/usr/local/bin/subnetwork.py stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
