# files/image_config/system-health/system-health.service
[Unit]
Description=SONiC system health monitor
Requires=database.service config-setup.service
After=database.service config-setup.service

[Service]
ExecStart=/usr/local/bin/healthd
Restart=always

[Install]
WantedBy=multi-user.target
