# files/image_config/udev/rules.d/99-rps.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="Ethernet*", RUN+="/usr/bin/rps.py %k"
