# platform/broadcom/sonic-platform-modules-nexthop/nh-4010/etc/systemd/system/system-health.service.d/override.conf
[Unit]
After=pddf-platform-init.service
